# config/quadlets/ollama.container
[Unit]
Description=Ollama LLM Service
After=local-fs.target

[Container]
Image=docker.io/ollama/ollama:latest
ContainerName=ollama
PublishPort=11434:11434
Volume=ollama-data:/root/.ollama
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=CUDA_VISIBLE_DEVICES=0
Device=nvidia.com/gpu=all

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
