# config/quadlets/ollama-webui.container
[Unit]
Description=Open WebUI (Ollama Frontend)
After=ollama.service
Requires=ollama.service

[Container]
Image=ghcr.io/open-webui/open-webui:main
ContainerName=ollama-webui
PublishPort=8080:8080
Volume=open-webui-data:/app/backend/data
Environment=OLLAMA_BASE_URL=http://host.containers.internal:11434

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
